# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; ----------------------------
; PHP-FPM
; ----------------------------
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/www.conf
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ----------------------------
; Nginx
; ----------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ----------------------------
; Laravel Scheduler
; ----------------------------
[program:laravel-scheduler]
command=/bin/bash -c "while true; do php /var/www/html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/html/storage/logs/scheduler-error.log
stderr_logfile_maxbytes=0

; ----------------------------
; Laravel Queue Worker
; ----------------------------
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log
stdout_logfile_maxbytes=0
stopwaitsecs=3600

; ----------------------------
; Laravel Setup (Migrations + Optimization)
; ----------------------------
[program:laravel-setup]
command=/bin/bash -c "\
until php /var/www/html/artisan migrate --force; do \
  echo 'Waiting for DB to be ready...'; sleep 5; \
done && php /var/www/html/artisan optimize:clear"
autostart=true
autorestart=false
stdout_logfile=/var/www/html/storage/logs/laravel-setup.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/html/storage/logs/laravel-setup-error.log
stderr_logfile_maxbytes=0
user=www-data
